# board/nftables/rootfs/etc/inittab
::sysinit:/etc/rc.local
::shutdown:/etc/rc.shutdown
::ctrlaltdel:/bin/true
