# board/system/rootfs/etc/inittab
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        : tty to run on, or empty for /dev/console
# runlevels : ignored
# action    : one of sysinit, respawn, askfirst, wait, and once
# process   : program to run
#
# Note: BusyBox init doesn't support runlevels.
#       The field is completely ignored.
#

# Start container system
::sysinit:/etc/init.d/rcS

# finish by calling rc.local, if it exists and is executable
::sysinit:/etc/rc.local

# Stuff to do before rebooting
::shutdown:/etc/init.d/rcK
